# bndl-stage/templates/default/logstash.init.erb
#!/usr/bin/env bash
#
# logstash
#
# chkconfig:   - 57 47
# description: logstash
# processname: logstash
# config:      <%= node[:elasticsearch][:path][:conf] %>/elasticsearch.yml

PIDFILE='/usr/local/logstash.pid'
LOGSTASH_USER =  <%= node[:logstash][:agent][:user] %>
CHECK_PID_RUNNING=$(ps ax | grep 'java' | grep -e "logstash.conf" | sed 's/^\s*\([0-9]*\)\s.*/\1/')

start() {
    if [ -f $PIDFILE ]; then
      # PIDFILE EXISTS -- ES RUNNING?
      echo -e "\033[31;1mPID file found in $PIDFILE, logstash already running?\033[0m"
      es_pid="$(cat $PIDFILE)"
      pid_running="$( ps ax | grep 'java' | grep logstash.conf )"

      if [ ! -z "$pid_running" ] ; then
        # EXIT IF ES IS ALREADY RUNNING
	      echo -e "\033[31;1mPID $es_pid still alive, already running...\033[0m"
	      return 1
      fi
    fi

    echo -e "\033[1mStarting logstash...\033[0m"
    touch $PIDFILE && chown <%= node[:logstash][:agent][:user] %> $PIDFILE
    su $LOGSTASH_USER -c "/usr/local/logstash/bin/logstash -f /usr/local/logstash/logstash.conf > '1' ? ' -d' : '' %> -p $PIDFILE"

    return $?
}

stop() {
    if [[ -f $PIDFILE ]]; then
      echo -n -e "\033[1mStopping logstash...\033[0m"

      # REMOVE PIDFILE AND EXIT IF PROCESS NOT RUNNING
      if [ ! $CHECK_PID_RUNNING ]; then
        echo -e "\033[1mPID file found, but no matching process running?\033[0m"
        echo    "Removing PID file..."
        su $LOGSTASH_USER -m -c "rm $PIDFILE"
        exit 0
      fi

      # KILL PROCESS
      pid=$(cat $PIDFILE)
      su $LOGSTASH_USER -m -c "kill $(cat $PIDFILE)"
      r=$?

      # Check for process
      timeout=0
      while [ $(ps -p $pid | wc -l ) -gt '1' ]; do
        echo -n '.'
        (( timeout ++))
        if [ $timeout -gt '30' ]; then return; fi
        sleep 1
      done;

      # Check for pidfile
      timeout=0
      while [ -f $PIDFILE  ]; do
        echo -n '.'
        (( timeout++ ))
        if [ $timeout -gt '15' ]; then return; fi
        sleep 1
      done; echo
      return $r
    else
      echo -e "\033[1mNo PID file found -- elasticsearch not running?\033[0m"
    fi
}

restart() {
    stop
    timeout=30
    while ps aux | grep 'java' | grep -e "logstash.conf"; do
      echo -n '.'
      (( timeout-- ))
      if [ $timeout -lt '1' ]; then return; fi
      sleep 1
    done;
    start
}


case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  restart)
        restart
        ;;
  *)
        echo $"Usage: $0 {start|stop|restart|}"
        exit 1
esac

exit $?
